# gnosis-backend.service
[Unit]
Description=Gnosis Trading System Backend API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/user/FINAL_GNOSIS
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 -m uvicorn web_api:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=append:/home/user/FINAL_GNOSIS/logs/backend.log
StandardError=append:/home/user/FINAL_GNOSIS/logs/backend.log

[Install]
WantedBy=multi-user.target
